# system-files/usr/local/etc/supervisor.d-deprecated/repo_setuppy.ini
[program:repo_setuppy]
command=python repo_setuppy.py
directory=/Users/russianidiot/git/andrewp-as-is/dotfiles/flask
autorestart=true
stopasgroup=true
killasgroup=true
stopsignal=INT
