# system-files/usr/local/etc/supervisor.d-deprecated/repo_travis_yml.ini
[program:repo_travis_yml]
command=python repo_travis_yml.py
directory=/Users/russianidiot/git/andrewp-as-is/dotfiles/flask
autorestart=true
stopasgroup=true
killasgroup=true
stopsignal=INT
