# system-files/usr/local/etc/supervisor.d-deprecated/repo_requirements.ini
[program:repo_requirements]
command=python repo_requirements.py
directory=/Users/russianidiot/git/andrewp-as-is/dotfiles/flask
autorestart=true
stopasgroup=true
killasgroup=true
stopsignal=INT
