# system-files/usr/local/etc/supervisor.d-deprecated/repo_readme.ini
[program:repo_readme]
command=python repo_readme.py
directory=/Users/russianidiot/git/andrewp-as-is/dotfiles/flask
autorestart=true
stopasgroup=true
killasgroup=true
stopsignal=INT

stderr_logfile=/Users/russianidiot/stderr_logfile.log
stdout_logfile=/Users/russianidiot/stdout_logfile.log
